# librenms.cron
# It's recommended not to run this cron anymore - please see librenms.nonroot.cron
# PLEASE DO NOT USE THIS FILE

33  */6   * * *   root    /opt/librenms/discovery.php -h all >> /dev/null 2>&1
*/5  *    * * *   root    /opt/librenms/discovery.php -h new >> /dev/null 2>&1
*/5  *    * * *   root    /opt/librenms/cronic /opt/librenms/poller-wrapper.py 16
15   0    * * *   root    /opt/librenms/daily.sh >> /dev/null 2>&1
*    *    * * *   root    /opt/librenms/alerts.php >> /dev/null 2>&1
