# poller-service.conf
# poller-service - SNMP polling service for LibreNMS

description "SNMP polling service for LibreNMS"
author "Clint Armstrong <clint@clintarmstrong.net>"

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
respawn

# Essentially lets upstart know the process will detach itself to the background
#expect fork

# # Run before process
# pre-start script
#     [ -d /var/run/myservice ] || mkdir -p /var/run/myservice
#     echo "Put bash code here"
# end script

chdir /opt/librenms
setuid librenms
setgid librenms

# Start the process
exec /opt/librenms/poller-service.py
