# dist/librenms-scheduler.timer
[Unit]
Description=LibreNMS Scheduler

[Timer]
OnCalendar=minutely
AccuracySec=1second

[Install]
WantedBy=timers.target
