# dist/librenms-scheduler.service
[Unit]
Description=LibreNMS Scheduler

[Service]
Type=oneshot
StandardOutput=null
StandardError=null
WorkingDirectory=/opt/librenms/
ExecStart=/usr/bin/env php artisan schedule:run --no-ansi --no-interaction
User=librenms
Group=librenms
